# config/supervisor/proxy.ini
[program:nginx]
command = /usr/sbin/nginx
exitcodes=0
autostart=true
autorestart=unexpected
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
